# ubuntu-development/s6/services.d/sshd/run
#!/command/with-contenv bash

 mkdir -p /run/sshd
 
exec 2>&1
exec /usr/sbin/sshd -D